# modules/wr/files/rq-dashboard.conf
description "Python RQ worker"
author      "Konrad Scherer"

start on runlevel [2345]
stop on runlevel [!2345]

setuid rq
chdir /home/rq/

exec /home/rq/.local/bin/rq-dashboard

