# modules/wr/files/rqworker-wr-rq.conf
description "Python RQ worker"
author      "Konrad Scherer"

start on runlevel [2345]
stop on runlevel [!2345]

setuid rq
chdir /home/rq/wr-jira-integration

exec /home/rq/.local/bin/rqworker wr-rq

